# supervisord.conf
[program:mailcatcher]
command=mailcatcher -f --smtp-port 2525 --smtp-ip 0.0.0.0 --http-ip 0.0.0.0
